# config/systemd/core.service.erb
[Unit]
Description=<%= fetch(:application) %> application server
After=network.target remote-fs.target nss-lookup.target
Wants=postgresql.service dneslov-sidekiq.service
Requires=redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
Environment=RAILS_ENV=<%= fetch(:rails_env) %>
WorkingDirectory=<%= current_path %>
ExecStart=/bin/bash -lc 'exec bundle exec puma'
StandardInput=null
StandardOutput=syslog
StandardError=syslog
TimeoutSec=30
SyslogIdentifier=dneslov

[Install]
WantedBy=multi-user.target
