# config/systemd/core.service.erb
[Unit]
Description=<%= fetch(:application) %> application core server
After=network.target remote-fs.target nss-lookup.target
Requires=redis.service postgresql.service dneslov_sidekiq.service

[Service]
Type=simple
User=_nginx
Group=_nginx
Environment=RAILS_ENV=<%= fetch(:rails_env) %>
WorkingDirectory=<%= current_path %>
ExecStart=/bin/bash -lc 'bundle exec rails s'
StandardInput=null
StandardOutput=syslog
StandardError=syslog
TimeoutSec=30
SyslogIdentifier=dneslov-core

[Install]
WantedBy=multi-user.target
