# config/systemd/web.service.erb
[Unit]
Description=<%= fetch(:application) %> application rails web-server
After=network.target remote-fs.target nss-lookup.target
Requires=postgresql.service redis.service dneslov_sidekiq.service

[Service]
Type=notify
User=www-data
Group=www-data
Environment=RAILS_ENV=<%= fetch(:rails_env) %>
Environment=PWD=<%= current_path %>
WorkingDirectory=<%= current_path %>
ExecStart=/bin/bash -lc 'exec bundle exec rails s'
StandardInput=null
StandardOutput=syslog
StandardError=syslog
TimeoutSec=10
SyslogIdentifier=dneslov-web

[Install]
WantedBy=multi-user.target
